# deploy/systemd/cdn-edge.service
[Unit]
Description=CDN Edge Server
Documentation=https://github.com/rzawadzk/cdn-edge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cdn-edge
Group=cdn-edge

# Environment-driven config (12-factor). Override in /etc/default/cdn-edge.
EnvironmentFile=-/etc/default/cdn-edge
ExecStart=/usr/local/bin/cdn-edge

# Restart policy.
Restart=always
RestartSec=2
StartLimitInterval=60
StartLimitBurst=5

# Graceful shutdown: send SIGTERM, give the process time to drain.
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Filesystem and process hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

# Allow writing to the cache dir (adjust path to match CDN_DISK_CACHE_DIR).
ReadWritePaths=/var/cache/cdn-edge

# Resource limits — tune to your hardware.
LimitNOFILE=1048576
LimitNPROC=65536

# Allow binding to privileged ports (80/443) without root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
